# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /webapps/django/uwsgi.ini
user=django

[program:nginx-app]
command = /usr/sbin/nginx
#user=django
