# puppet/files/gunicorn-supervisord.conf
[program:gunicorn]
command=gunicorn mysite.wsgi:application -c /webapps/mysite/gunicorn.conf.py
directory=/webapps/mysite
user=mysite
autostart=true
autorestart=true
redirect_stderr=true
