# puppet/files/gunicorn-supervisord.conf
[program:gunicorn]
command=gunicorn mysite.wsgi:application -c /webapps/django/gunicorn.conf.py
directory=/webapps/django
user=django
autostart=true
autorestart=true
redirect_stderr=true
